# var/service/pulseaudio/run
#!/data/data/com.termux/files/usr/bin/sh
exec pulseaudio --daemonize=no --log-target=stderr 2>&1
